# 3.0.1.0-php7.1/scripts/supervisord-zopencart.conf
[program:opencart-setup]
command=sh -x /setup-opencart.sh
startsecs=1
autorestart=false
startretries=1